# res/backup-tracking-db.service
[Unit]
Description=Backup tracking.db

[Service]
ExecStart=/usr/bin/cp /path/to/tracking.db /path/to/backup/dir/tracking.db
